# supervisor/beat.conf
[program:beat]
command=/usr/local/bin/sentry --config=/srv run beat
directory=/srv
user=sentry
autostart=true
autorestart=true
