# supervisor/cron.conf
[program:beat]
command=/usr/local/bin/sentry --config=/srv run cron
directory=/srv
user=sentry
autostart=true
autorestart=true
