# supervisor/redis.conf
[program:redis]
command=/usr/bin/redis-server /srv/redis.conf
directory=/srv
user=sentry
autostart=true
autorestart=true
