# supervisor/sentry.conf
[program:sentry]
directory=/srv
command=/usr/local/bin/sentry --config=/srv/sentry.conf.py start
user=sentry
autostart=true
autorestart=true
