# supervisor/worker.conf
[program:worker]
command=/usr/local/bin/sentry --config=/srv/sentry.conf.py celery worker -B
directory=/srv
user=sentry
autostart=false
autorestart=true
