# supervisor/worker.conf
[program:worker]
command=/usr/local/bin/sentry --config=/srv celery worker
directory=/srv
user=sentry
autostart=true
autorestart=true
